# docker/phpldapadmin/copy/etc/supervisor/conf.d/phpldap.conf
[program:nginx]
command=/usr/sbin/nginx

[program:php]
command=/usr/sbin/php-fpm7.0 -F
