# docker/file_router/copy/etc/supervisor/conf.d/samba.conf
[program:samba]
command=/usr/sbin/smbd -D
