# docker/file_router/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:rpcbind]
command=rpcbind

#[program:mountd]
#command=rpc.mountd

#[program:nfsd]
#command=rpc.nfsd
