# docker/ldap/copy/etc/supervisor.d/ldap.ini
[program:slapd]
command=slapd -d 0 -h "ldap:/// ldaps:///"
