# docker/postgres/copy/etc/supervisor.d/postgres.conf
[program:postgres]
command=/usr/lib/postgresql/9.5/bin/postgres -D /var/lib/postgresql/data
user=postgres
