# docker/postgres/copy/etc/supervisor.d/postgres.ini
[program:postgres]
command=postgres -D /var/lib/postgresql/data/
user=postgres
