# docker/zabbix/copy/etc/supervisor.d/zabbix.ini
[program:nginx]
command=/usr/sbin/nginx

[program:php]
command=/usr/bin/php-fpm --nodaemonize

[program:zabbix]
command=/usr/sbin/zabbix_server -f
user=zabbix
