# docker/analysis/copy/etc/supervisor/conf.d/jupyterhub.conf
[program:jupyterhub]
command = /core/software/conda/bin/jupyterhub -f /etc/jupyterhub/jupyterhub_config.py
environmet = PATH="/vcore/software/conda/bin:%(ENV_PATH)s"
