# docker/analysis/copy/etc/supervisor/conf.d/jupyterhub.conf
[program:jupyterhub]
command = /vcore/software/conda/bin/jupyterhub -f /etc/jupyterhub/jupyterhub_config.py
environment = PATH="/vcore/software/conda/bin:%(ENV_PATH)s"
